# src/souzu/systemd/res/souzu-update.service.template
[Unit]
Description=Souzu self-updater
After=network-online.target

[Service]
Type=oneshot
ExecStart={souzu_path} update --restart
Restart=no
