# src/souzu/systemd/res/souzu.service.template
[Unit]
Description=Souzu Slack notifications for Bambu printers
After=network-online.target

[Service]
Type=simple
ExecStart={souzu_path} monitor
Restart=always
RestartSec=10

[Install]
WantedBy=default.target